# Customize/luci-app-dockerman/root/etc/uci-defaults/luci-app-dockerman
#!/bin/sh

/etc/init.d/dockerd disable

exit 0
